# app.service
[Unit]
Description=App Service (Backend + Frontend)
After=network.target

[Service]
Type=forking
User=root
WorkingDirectory=/root/video2sop
ExecStart=/root/video2sop/start_services_persistent.sh
ExecStop=/root/video2sop/stop_services.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
